# deploy/promptrelay-worker.service
[Unit]
Description=PromptRelay Worker
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/promptrelay-worker
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=3
Environment=NODE_ENV=production

# Optional env override examples (primary config is src/config.ts)
# Environment=PORT=8091
# Environment=GOOGLE_CLOUD_PROJECT=your-project-id

[Install]
WantedBy=multi-user.target
